# deploy/mcp-gsc-public.service
[Unit]
Description=MCP Google Search Console Public Server
After=network.target

[Service]
Type=simple
User=yusof
WorkingDirectory=/home/yusof/mcp/MCP-Public
EnvironmentFile=/home/yusof/mcp/MCP-Public/.env
ExecStart=/home/yusof/mcp/MCP-Public/.venv/bin/python gsc/run.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
